# target/linux/ramips/mt7620/base-files/etc/uci-defaults/05_fix-compat-version
. /lib/functions.sh

case "$(board_name)" in
	ravpower,rp-wd03)
		uci set system.@system[0].compat_version="2.0"
		uci commit system
		;;
esac
uci set wireless.default_radio1.disabled='0'
uci set wireless.default_radio1.ssid='Op2203_2.4G'
uci set wireless.default_radio1.key='12345678'
uci set wireless.default_radio1.encryption='psk-mixed'
uci set wireless.default_radio0.disabled='0'
uci set wireless.default_radio0.ssid='Op2203'
uci set wireless.default_radio0.key='12345678'
uci set wireless.default_radio0.encryption='psk-mixed'
uci set wireless.radio1.disabled='0'
uci set wireless.radio0.disabled='0'
uci commit wireless
uci set network.lan.ipaddr='192.168.6.1'
uci commit network


exit 0
